# salt/srv/salt/base/cron/files/crontab
*/1 * * * * echo "Hello from crontab file" >> /tmp/cronhello
